# cicd/tnik.service
[Unit]
Description=tnik

[Service]
ExecStart=/opt/tnik/server.sh

[Install]
WantedBy=multi-user.target